# deploy/services/novel-api.service
[Unit]
Description=Inherited Cloud — FastAPI Query Server
After=network.target

[Service]
Type=simple
User=novel
WorkingDirectory=/opt/inherited-cloud
EnvironmentFile=/opt/inherited-cloud/.env
ExecStart=/opt/inherited-cloud/.venv/bin/uvicorn api.server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
